# systemd/package-export.service
[Unit]
Description=Export installed package list to dotfiles repo

[Service]
Type=oneshot
WorkingDirectory=%h/dotfiles
ExecStart=%h/dotfiles/scripts/package-export.sh
