# systemd/repo-update-check.service
[Unit]
Description=Check local git repos for remote updates

[Service]
Type=oneshot
WorkingDirectory=%h
ExecStart=%h/dotfiles/scripts/repo-update-check.sh
