# systemd/system-update.service
[Unit]
Description=System package updater

[Service]
Type=oneshot
WorkingDirectory=%h/dotfiles
ExecStart=%h/dotfiles/scripts/system-update.sh
